# custom-kafka-exporter.service.template
[Unit]
Description={{ ServiceDescription }}
After=network.target

[Service]
Type=simple
User={{ User }}
ExecStart={{ Command }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
